# docker/pipeline-back/crontabs/ncanda
MAILTO=jim.klo@sri.com
0 21 * * *      /sibis-software/ncanda-data-integration/scripts/crond/back-nightly
