# docker/pipeline-back/crontabs/root
1 0 * * 0	/sibis-software/ncanda-data-integration/scripts/transmit/transmit_all
