# docker/pipeline-front/crontabs/ncanda
MAILTO=ncanda-admin@sri.com
0 * * * *      /sibis-software/ncanda-data-integration/scripts/crond/front-hourly
