# docker/pipeline-front/crontabs/ncanda
0 * * * *      /sibis-software/ncanda-data-integration/scripts/crond/front-hourly
